# src/skel/etc/systemd/system/fillmetadata.service
[Unit]
Description=Fillmetadata script as part of Small Files Plugin

[Service]
ExecStart=/usr/bin/python3 /usr/local/bin/fillmetadata.py
