# src/skel/etc/systemd/system/pack-files.service
[Unit]
Description=Pack-files script as part of Small Files Plugin

[Service]
ExecStart=/usr/bin/python3 /usr/local/bin/pack-files.py