# src/skel/etc/systemd/system/writebfids.service
[Unit]
Description=Writebfids script as part of Small Files Plugin

[Service]
ExecStart=/usr/bin/python3 /usr/local/bin/writebfids.py